# debian/libsites.cron.daily
#!/bin/bash

cd /srv/libsites
perl -Ilib -Ilocal/lib/perl5 bin/update all

