# debian/libsites.cron.daily
#!/bin/bash

cd /srv/libsites
make update

